# 50-parportser.rules
SUBSYSTEMS=="usb", DRIVERS=="pl2303", RUN+="/usr/src/parportser/load.sh /dev/ttyUSB0"
